# backend/labs/base-novnc/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:xvfb]
# Crea pantalla virtual en :1 con resolución 1280x720
command=Xvfb :1 -screen 0 1280x720x24
autorestart=true
priority=1                          ; arranca primero

[program:xfce4]
# Lanza escritorio Xfce4 en la pantalla virtual
command=startxfce4
environment=DISPLAY=":1"
autorestart=true
priority=2                          ; espera a que Xvfb esté listo

[program:x11vnc]
# Servidor VNC que lee la pantalla virtual y la expone en el puerto 5900
command=x11vnc -display :1 -nopw -listen localhost -xkb -forever
autorestart=true
priority=3

[program:novnc]
# websockify hace de puente: convierte el WebSocket del navegador a VNC
# --web apunta a los ficheros estáticos de noVNC (la UI del navegador)
command=websockify --web /usr/share/novnc 6080 localhost:5900
autorestart=true
priority=4                          ; arranca último, cuando VNC ya está listo