# labs/base-novnc/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:xvfb]
command=Xvfb :1 -screen 0 1280x720x24
autorestart=true
startsecs=2
priority=1

[program:xfce4]
# dbus-run-session provides the session D-Bus daemon that xfce4 components require.
# Without it xfwm4/xfce4-panel/xfdesktop fail silently and the desktop is blank.
command=dbus-run-session startxfce4
environment=DISPLAY=":1"
autorestart=true
startsecs=3
startretries=5
priority=2

[program:x11vnc]
command=x11vnc -display :1 -nopw -listen localhost -xkb -forever
autorestart=true
startsecs=1
startretries=10
priority=3

[program:novnc]
# websockify bridges the browser WebSocket to the VNC port.
# Starts last so x11vnc is already listening on 5900.
command=websockify --web /usr/share/novnc 6080 localhost:5900
autorestart=true
startsecs=1
startretries=5
priority=4
